# templates/apache.service.j2
[Unit]
Description=Red Hat JWS Apache HTTPD
After=syslog.target network.target

[Service]
Type=forking
WorkingDirectory={{ jws_apache.home }}/sbin
Environment="LD_LIBRARY_PATH={{ jws_apache.home }}/lib"

ExecStart={{ jws_apache.home }}/sbin/apachectl start -f {{ jboss_jws_instance_dir }}/conf/httpd.conf
ExecStop={{ jws_apache.home }}/sbin/apachectl stop -f {{ jboss_jws_instance_dir }}/conf/httpd.conf

User=root
Group=root
#User={{ jws_apache.user }}
#Group={{ jws_apache.group }}

[Install]
WantedBy=multi-user.target
